# databases/influxdb/files/influxdb.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: influxdb
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable influxdb:
# influxdb_enabe="YES"
#
# influxdb_enable (bool):	Set to YES to enable influxdb
#				Default: NO
# influxdb_conf (str):		influxdb configuration file
#				Default: ${PREFIX}/etc/influxdb.conf
# influxdb_flags (str):		Extra flags passed to influxdb


. /etc/rc.subr

name="influxdb"
rcvar=influxdb_enable

command="%%PREFIX%%/bin/influxdb"
command_args='> /dev/null 2>&1 &'

load_rc_config $name

[ -z "$influxdb_enable" ]	&& influxdb_enable="NO"
[ -z "$influxdb_conf" ]		&& influxdb_conf="%%PREFIX%%/etc/influxdb.conf"
[ -z "$influxdb_flags" ]	&& influxdb_flags=""
[ -n "$influxdb_conf" ]		&& influxdb_flags="$influxdb_flags -config=$influxdb_conf"

run_rc_command "$1"
